# tests/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:httpd]
command = /usr/sbin/httpd -DFOREGROUND
stdout_logfile = /var/log/supervisor/httpd.log
user = root
autostart = true

[program:php-fpm]
command = /usr/sbin/php-fpm --nodaemonize
stdout_logfile = /var/log/supervisor/php-fpm.log
user = root
autostart = true

[program:mysqld]
command = /usr/bin/pidproxy /var/run/mariadb/mariadb.pid /usr/bin/mysqld_safe
stdout_logfile = /var/log/supervisor/mysqld.log
redirect_stderr = true
user = mysql
autostart = true

[program:salt-master]
command = /usr/bin/salt-master
stdout_logfile = /var/log/supervisor/salt-master.log
redirect_stderr = true
user = root
autostart = true
